# php-fpm/supervizord.cfg
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)
username=www-data
password=

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket
username=www-data
password=

[program:gp_tpl_worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/app/artisan queue:work --tries=3
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=8
redirect_stderr=true
stdout_logfile=/var/log/worker.log
stopwaitsecs=3600

[program:gp_tpl_scheduler]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/app/artisan schedule:work
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=8
redirect_stderr=true
stdout_logfile=/var/log/scheduler.log
stopwaitsecs=3600
